# video-forge-bot.service
[Unit]
Description=Video Forge Telegram Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/opt/video-forge-bot
Environment=PATH=/opt/video-forge-bot/venv/bin
ExecStart=/opt/video-forge-bot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=video-forge-bot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/video-forge-bot/data
ReadWritePaths=/opt/video-forge-bot/user_credits.json

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
